# etc/init/hk_go.conf
#!upstart
description "Startup script for h&k site (go version)"
author      "hector@hectorcorrea.com"
# source:   http://dailyjs.com/2011/03/07/node-deployment/
# location: /etc/init/hk_go.conf

respawn
start on runlevel [23]

script
  export DB_USER="***"
  export DB_PASSWORD="***"
  export DB_NAME="***"

  export BLOG_USR="***"
  export BLOG_PASS="***"
  export BLOG_SALT="***"

  cd /root/hk_go
  exec ./hk -address localhost:3003 >> /var/log/hk_go.log 2>&1
end script
